# meta-dgam-pr/recipes-app/kubesolo/files/kubesolo.service
[Unit]
Description=Kubesolo Kubernetes Cluster
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple

# Set kubeconfig location
Environment="KUBECONFIG=/var/lib/kubesolo/pki/admin/admin.kubeconfig"

# Wrapper script handles config polling and then execs kubesolo
ExecStart=/usr/bin/kubesolo-start.sh

# Only restart if kubesolo itself crashes after a successful start
Restart=on-failure
RestartSec=60s

[Install]
WantedBy=multi-user.target
